# build/rc.local
sudo /home/pi/syncraft-scripts/scripts/main.sh