# build/rc.local
#!/bin/bash -e

export DISPLAY=:0.0

sudo python3 /home/pi/syncraft-scripts/scripts/main.py