# build/rc.local
#!/bin/bash -e

export DISPLAY=:0.0

sudo python3 ~/syncraft-scripts/scripts/startup.py