# build/syncraft-scripts.service
[Unit]
Description=Scripts to Syncraft 3D printers
Before=KlipperScreen.service

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/python3 /home/pi/syncraft-scripts/scripts/main.py
User=pi